# sv/collector/run
#!/bin/bash

sv start cassandra || exit 1

source /etc/envvars
env

cd /zipkin

exec 2>&1
exec ./bin/collector cassandra
