# sv/query/run
#!/bin/bash

sv start mysql || exit 1

source /etc/envvars
export MYSQL_USER=root 
env

cd /zipkin

exec 2>&1
exec ./bin/query mysql
